# deploy/bolshoi-bot.service
[Unit]
Description=Bolshoi Theatre Bot for Telegram Messenger
After=network.target

[Service]
Type=simple
User=bolshoi-bot
WorkingDirectory=/home/bolshoi-bot/chip_bolshoi_telegram
EnvironmentFile=/home/bolshoi-bot/chip_bolshoi_telegram/.env
ExecStart=/home/bolshoi-bot/chip_bolshoi_telegram/venv/bin/python -m bot.main
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/bolshoi-bot/bot.log
StandardError=append:/var/log/bolshoi-bot/bot.log

[Install]
WantedBy=multi-user.target
